# sebastian.service
[Unit]
Description=SebasTian Assistant
After=network.target

[Service]
User=your_username
WorkingDirectory=/path/to/sebastian
Environment="PATH=/path/to/your/venv/bin"
ExecStart=/path/to/your/venv/bin/uvicorn src.main:app --host localhost --port 8000

[Install]
WantedBy=multi-user.target 